# templates/default/nexus.erb
#!/bin/sh

#
# chkconfig: - 85 15
#

export RUN_AS_USER=nexus
/usr/local/bin/nexus "$@"
exit $?